# src/docker/supervisor-memcached.ini
[program:memcached]
command=/usr/bin/memcached -u memcached -v -m 64 -p 11211 -c 1024
autostart=true
autorestart=true
redirect_stderr=true
user=root

