# src/docker/supervisor-apache2.ini
[program:apache2]
command=httpd -D "FOREGROUND" -k start
numprocs=1
autostart=true
autorestart=true
startretries=1
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/apache2.log
stderr_logfile=/var/log/apache2.err.log
user=root
